# infra/supervisor/blabprezent.conf
[program:blabprezent]
command=/usr/bin/python3 /home/bp/WebsiteParsingBot.py
directory=/home/bp
user=bp

