# infra/supervisor/blabprezent.conf
[program:blabprezent]
environment=PATH="/home/bp/bpenv/bin"
command=/home/bp/bpenv/bin/python /home/bp/src/bp.py bot
directory=/home/bp/src
user=bp
